# cronitor/crontab
* * * * * curl -m 5 --retry 5 https://cronitor.link/p/RANDOMSTRING/dashubnt-VM-heartbeat
