# cronitor/crontab
* * * * * curl -m 5 --retry 5 https://cronitor.link/p/RANDOMSTRING/dashubnt-VM-heartbeat
0 12 * * * cronitor exec gaVpgM /root/disk-bot/full-disk-bot.sh
